# scripts/kaseki-api.service
[Unit]
Description=Kaseki Agent API Service
Documentation=https://github.com/CyanAutomation/kaseki-agent
After=network.target docker.service

[Service]
Type=simple
User=nobody
Group=nogroup
WorkingDirectory=/agents/kaseki-template
EnvironmentFile=-/etc/kaseki-api/kaseki-api.env

# Security hardening (match Docker settings)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/agents/kaseki-results /var/log/kaseki-api /tmp
UMask=0077

# Restart behavior
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kaseki-api

# Timeout and graceful shutdown
TimeoutStopSec=30
KillSignal=SIGTERM

# NOTE: Actual command depends on deployment method:
# Option 1: Run as Node.js process (on host)
# ExecStart=/usr/bin/node /agents/kaseki-template/dist/kaseki-api-service.js

# Option 2: Run in Docker container
ExecStartPre=/usr/bin/bash -lc 'printf "KASEKI_API_IMAGE=%s\n" "${KASEKI_API_IMAGE:-kaseki-agent:node24-local}"; docker image inspect "${KASEKI_API_IMAGE:-kaseki-agent:node24-local}" --format "ImageId={{.Id}} RepoDigests={{join .RepoDigests \",\"}}"'
ExecStartPre=/usr/bin/docker run --rm --entrypoint /bin/bash ${KASEKI_API_IMAGE:-kaseki-agent:node24-local} -lc 'test -f /app/dist/kaseki-api-service.js && ls -l /app/dist/kaseki-api-service.js'
ExecStart=/usr/bin/docker run --rm \
  --name kaseki-api-container \
  -p 8080:8080 \
  -v /agents:/agents:rw \
  -v /agents/kaseki-results:/agents/kaseki-results:rw \
  -v /var/log/kaseki-api:/var/log/kaseki-api:rw \
  -v /var/run/docker.sock:/var/run/docker.sock \
  --env-file /etc/kaseki-api/kaseki-api.env \
  --user 0:0 \
  --cap-drop ALL \
  --security-opt no-new-privileges:true \
  --read-only \
  --entrypoint /bin/bash \
  ${KASEKI_API_IMAGE:-kaseki-agent:node24-local} \
  -lc 'set -euo pipefail; test -f /app/dist/kaseki-api-service.js; ls -l /app/dist/kaseki-api-service.js; exec node /app/dist/kaseki-api-service.js'

ExecStop=/usr/bin/docker stop -t 30 kaseki-api-container

[Install]
WantedBy=multi-user.target
